# resources/safeshutdown/SafeShutdown.service
[Unit]
After=network.service
[Service]
ExecStart=/usr/sbin/safeshutdown-wrapped.sh
[Install]
WantedBy=default.target
